# bootstrap/bootup.service
[Unit]
Description=Bootup Script

[Service]
ExecStart=/usr/local/bin/bootup.sh

[Install]
WantedBy=default.target
